# src/launcher/tv-launcher-bridge.service
[Unit]
Description=TV Launcher Bridge Server
After=network.target graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/launcher/server.py
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
